# root/etc/services.d/teamspeak_helper/run
#!/usr/bin/with-contenv sh
while ! id -u ts >/dev/null 2>&1; do
        sleep 1s
done

exec s6-setuidgid ts /teamspeak/sh/helper.sh