# root/etc/services.d/teamspeak_helper/run
#!/usr/bin/with-contenv sh
sleep 10s
exec s6-setuidgid ts /teamspeak/sh/helper.sh