# root/etc/services.d/teamspeak/finish
#!/usr/bin/with-contenv sh
if [ -e "/teamspeak/ts3server.ini" ]
then
  mv /teamspeak/ts3server.ini /teamspeak/save/ts3server.ini
else
  s6-svscanctl -t /var/run/s6/services
fi