# root/etc/services.d/teamspeak/finish
#!/bin/sh
if [ -e "/teamspeak/init" ]
then
  rm -R /teamspeak/init
fi
s6-svscanctl -t /var/run/s6/services
