# root/etc/services.d/teamspeak/finish
#!/usr/bin/with-contenv sh
if [ -e /teamspeak/init ]
then
  rm -R /teamspeak/init
else
  s6-svscanctl -t /var/run/s6/services
fi