# root/etc/services.d/teamspeak/run
#!/usr/bin/with-contenv sh
sh /teamspeak/sh/startup.sh
s6-echo "Starting TeamSpeak 3 Server.."
exec s6-setuidgid ts /teamspeak/ts3server_minimal_runscript.sh