# root/etc/services.d/teamspeak/run
#!/usr/bin/with-contenv sh
sh /teamspeak/sh/startup.sh

if [ -e "/teamspeak/ts3server_minimal_runscript.sh" ]
then
        exec s6-setuidgid ts /teamspeak/ts3server_minimal_runscript.sh    
fi